# test-appliance/files/etc/systemd/system/gce-finalize.timer
[Unit]
Description=GCE finalization timer

[Timer]
#OnBootSec=23h 45m
OnBootSec=5d
Unit=gce-finalize.service

[Install]
WantedBy=timers.target
